# roles/nodejs/templates/etc/profile.d/npm.sh.j2
# {{ ansible_managed }}

export PATH={{ npm_config_directory }}/bin:$PATH
export npm_config_prefix={{ npm_config_directory }}
export NODE_PATH=$NODE_PATH:{{ npm_config_directory }}/lib/node_modules
